# node_exporter.service
[Unit]
Description=Service which runs Node Exporter for Prometheus scraping
 
[Service]
User=ergonized
ExecStart=/usr/bin/node_exporter
 
[Install]
WantedBy=multi-user.target
